# app/docker/supervisord.conf
; bid-app supervisord 配置
; 参考 IMPLEMENTATION_SPEC §17.2
; 注意:alembic upgrade 已由 entrypoint.sh 同步执行(D-O),supervisord 只负责长期进程

[supervisord]
nodaemon=true
logfile=/var/log/supervisord.log
pidfile=/var/run/supervisord.pid

[program:uvicorn]
command=/app/backend/.venv/bin/uvicorn bid_app.main:app --host 0.0.0.0 --port 12123 --no-access-log
directory=/app/backend
autostart=true
autorestart=true
priority=20
stopsignal=TERM
stopwaitsecs=10
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0

[program:arq-worker]
command=/app/backend/.venv/bin/arq bid_app.worker.settings.WorkerSettings
directory=/app/backend
autostart=true
autorestart=true
priority=20
stopsignal=TERM
stopwaitsecs=30
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0

[program:cron]
command=/usr/sbin/cron -f
autostart=true
autorestart=true
priority=30
stdout_logfile=/dev/fd/1
stderr_logfile=/dev/fd/2
